# cron/weekday.cron
# ┌─min ─hour ─dom ─mon ─dow
0   13   *    *    1-5  python /app/hearing_bot.py check_status
                        python /app/hearing_bot.py update && \
                        