# cron/weekly.cron
0   14   *    *    1    python /app/hearing_bot.py post_upcoming
